# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create minicluster user if it doesn't exist
        if ! getent group minicluster > /dev/null 2>&1; then
            groupadd --system minicluster
        fi
        if ! getent passwd minicluster > /dev/null 2>&1; then
            useradd --system --gid minicluster --home-dir /opt/minicluster \
                --shell /usr/sbin/nologin minicluster
        fi

        # Create data directories
        mkdir -p /opt/minicluster/api/data
        mkdir -p /var/log/minicluster

        # Set permissions
        chown -R minicluster:minicluster /opt/minicluster
        chown -R minicluster:minicluster /var/log/minicluster
        chmod 750 /opt/minicluster/api/data

        # Link configuration
        if [ ! -L /opt/minicluster/api/appsettings.json ]; then
            ln -sf /etc/minicluster/appsettings.json /opt/minicluster/api/appsettings.json
        fi

        # Enable and start service
        systemctl daemon-reload
        systemctl enable minicluster.service
        
        echo ""
        echo "========================================"
        echo "  MiniCluster installed successfully!"
        echo "========================================"
        echo ""
        echo "Start the service:"
        echo "  sudo systemctl start minicluster"
        echo ""
        echo "Access the dashboard at:"
        echo "  http://localhost:5147"
        echo ""
        echo "Configuration file:"
        echo "  /etc/minicluster/appsettings.json"
        echo ""
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
