# packaging/debian-go/minicluster.service
[Unit]
Description=MiniCluster Control Center API
Documentation=https://github.com/innovatek/minicluster
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/minicluster

ExecStart=/usr/bin/minicluster-api

Restart=always
RestartSec=10

# Environment — override in /etc/minicluster/config.yaml or via environment
Environment=MINICLUSTER_DATA_DIR=/var/lib/minicluster
Environment=MINICLUSTER_CONFIG_FILE=/etc/minicluster/config.yaml

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=minicluster

# Resource limits
LimitNOFILE=65536
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
